# genie-internal/jenie-genie.service
[Unit]
Description=Jenie Genie - Internal
After=network.target

[Service]
# Folder where your code and .env.internal live
WorkingDirectory=/opt/jenie-genie

# Use the venv's Python and your entry script
ExecStart=/opt/jenie-genie/venv/bin/python server_internal.py

# Restart if it crashes
Restart=on-failure
RestartSec=5

# (Optional) If your env loader wants a specific env file name:
# Environment=GENIE_ENV_FILE=.env.internal

[Install]
WantedBy=default.target
